# archiso/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear --noissue %I $TERM
TTYVTDisallocate=no
